# telegram-bot.service
[Unit]
Description=Telegram Bot - Qarz Eslatmalari
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=jamshid
Group=jamshid
WorkingDirectory=/var/www/jamshid
Environment="PATH=/var/www/jamshid/venv/bin"
ExecStart=/var/www/jamshid/venv/bin/python3 /var/www/jamshid/run_telegram_bot.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-bot

[Install]
WantedBy=multi-user.target
